# sqlrelay_frequency.service
[Unit]
Description=Delete Recent Queries every 500 seconds

[Service]
ExecStart=/usr/local/firstworks/deleterecent.sh
